# py/test/pytests/touchscreen_calibration/sysfs_service.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Responsible for running sysfs_server

description   "Run sysfs_server"
author        "chromium-os-dev@chromium.org"

start on started shill

script
  # Wait for network to be completely ready.
  sleep 10

  PROJ_DIR="/var/tmp/ts"
  MACHINE_IP=`ifconfig eth0 | egrep "inet .+netmask .+broadcast" |\
              awk '{print $2}'`
  logger -t "${UPSTART_JOB}" "MACHINE_IP: ${MACHINE_IP}"
  exec sudo python ${PROJ_DIR}/sysfs_server.py ${MACHINE_IP}:8000
end script
