# py/test/pytests/touchscreen_calibration/touch_sensors_service.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Responsible for running sensors_server

description   "Run sensors_server"
author        "chromium-os-dev@chromium.org"

start on started shill

script
  # Wait for network to be completely ready.
  sleep 10

  PROJ_DIR="/var/tmp/ts"
  BOARD_PATH="$PROJ_DIR/boards/board"
  MACHINE_IP=`ifconfig eth0 | egrep "inet .+netmask .+broadcast" |\
              awk '{print $2}'`
  # Must specify the target board if this program is run on a beagle bone.
  BOARD=`test -f "$BOARD_PATH" && cat "$BOARD_PATH" ||\
         cat /etc/lsb-release | grep CHROMEOS_RELEASE_BOARD |\
         awk -F "=" '{print $2}'`
  logger -t "${UPSTART_JOB}" "MACHINE_IP: ${MACHINE_IP}"
  logger -t "${UPSTART_JOB}" "BOARD: ${BOARD}"
  sudo mount -o remount,rw /
  exec sudo python ${PROJ_DIR}/sensors_server.py ${MACHINE_IP}:8000 ${BOARD}
end script
